# configs/com.user.logwriter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.logwriter</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/log_writer.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/opt/app/service.log</string>
    <key>StandardErrorPath</key>
    <string>/opt/app/service.error.log</string>
</dict>
</plist>
